# assets/configs/websockify.conf
[program:websockify]
command=/system/assets/configs/websockify.sh
autorestart=true
